# etc/udev/rules.d/50-udev.rules
# There are a number of modifiers that are allowed to be used in some of the
# fields.  See the udev man page for a full description of them.
#
# See the udev.rules.examples file for more examples of how to create rules
#

# block devices
BUS=="ide", KERNEL=="*", SYSFS{removable}=="1", NAME="cdrom", SYMLINK="cdroms/%k"
BUS=="scsi", KERNEL=="sr0", SYSFS{removable}=="1", NAME="cdrom", SYMLINK="cdroms/%n"
BUS=="scsi", KERNEL=="sr[1-9]", SYSFS{removable}=="1", NAME="cdrom%n", SYMLINK="cdroms/%n"
BUS=="scsi", KERNEL=="*", NAME="%k"
BUS=="ide", KERNEL=="*", NAME="%k"

# create a symlink named after the device map name
# note devmap_name comes with extras/multipath
#KERNEL=="dm-[0-9]*", PROGRAM="/sbin/devmap_name %M %m", NAME="%k", SYMLINK="%c"

# DRI devices always go into a subdirectory (as per the LSB spec)
KERNEL=="card*", NAME="dri/card%n"

# alsa devices
KERNEL=="controlC[0-9]*", NAME="snd/%k"
KERNEL=="hw[CD0-9]*",     NAME="snd/%k"
KERNEL=="pcm[CD0-9cp]*",  NAME="snd/%k"
KERNEL=="midiC[D0-9]*",   NAME="snd/%k"
KERNEL=="timer",          NAME="snd/%k"
KERNEL=="seq",            NAME="snd/%k"

# USB devices
KERNEL=="hiddev*",	NAME="usb/%k"
KERNEL=="auer*",		NAME="usb/%k"
KERNEL=="legousbtower*",	NAME="usb/%k"
KERNEL=="dabusb*",	NAME="usb/%k"
BUS=="usb", KERNEL=="lp[0-9]*",	NAME="usb/%k",		SYMLINK="printers/usb%n"

# CAPI devices
KERNEL=="capi",		NAME="capi20", SYMLINK="isdn/capi20"
KERNEL=="capi*",		NAME="capi/%n"

KERNEL=="mISDN",		NAME="mISDN", SYMLINK="isdn/mISDN"

# Network devices
KERNEL=="tun",		NAME="net/%k"

KERNEL=="lp[0-9]",   NAME="%k",        SYMLINK="printers/%n";

# md block devices
KERNEL=="md[0-9]*", NAME="md/%n"

# floppy devices
KERNEL=="fd[0-9]*", NAME="%k",			SYMLINK="floppy/%n"

# tty devices
KERNEL=="tty[0-9]*",    NAME="%k"		SYMLINK="vc/%n"

#GSM Modems
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0408",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="04a5",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="05c6",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0af0",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0c88",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0eab",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="12d1",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1410",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1614",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="16d5",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="16d8",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1726",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="19d2",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1a8d",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1bc7",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="1da5",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="413c",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0421",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"
KERNEL=="ttyACM[0-9]*", SYSFS{idVendor}=="0e8d",   NAME="ttyGSM%n",   SYMLINK="tts/GSM%n"

#Sentech My Wireless
KERNEL=="ttyUSB[0-9]*", SYSFS{idVendor}=="0bc3",   NAME="ttyIPW%n",   SYMLINK="tts/IPW%n"

KERNEL=="ttyUSB[0-9]*", NAME="%k",              SYMLINK="tts/USB%n"

# vc devices
KERNEL=="vcs",        NAME="%k",		SYMLINK="vcc/0"
KERNEL=="vcs[0-9]*",  NAME="%k",		SYMLINK="vcc/%n"
KERNEL=="vcsa",       NAME="%k0",		SYMLINK="vcc/a0"
KERNEL=="vcsa[0-9]*", NAME="%k",		SYMLINK="vcc/a%n"

# v4l devices
KERNEL=="video[0-9]*", NAME="%k",	SYMLINK="v4l/video%n"
KERNEL=="radio[0-9]*", NAME="%k",	SYMLINK="v4l/radio%n"
KERNEL=="vbi[0-9]*",   NAME="%k",	SYMLINK="v4l/vbi%n"
KERNEL=="vtx[0-9]*",   NAME="%k",	SYMLINK="v4l/vtx%n"

# dm devices (ignore them)
KERNEL=="dm-[0-9]*", NAME=""

# i2c devices
KERNEL=="i2c-[0-9]*", NAME="%k",		SYMLINK="i2c/%n"

# loop devices
KERNEL=="loop[0-9]*", NAME="%k", SYMLINK="loop/%n"

# ramdisks
KERNEL=="ram[0-9]*", NAME="%k", SYMLINK="rd/%n"

# framebuffer devices
KERNEL=="fb[0-9]*", NAME="%k", SYMLINK="fb/%n"

# misc
KERNEL=="rtc",     NAME="%k", SYMLINK="misc/%k"
KERNEL=="psaux",   NAME="%k", SYMLINK="misc/%k"
KERNEL=="agpgart", NAME="%k", SYMLINK="misc/%k"
KERNEL=="rtc",     NAME="%k", SYMLINK="misc/%k"
KERNEL=="psaux",   NAME="%k", SYMLINK="misc/%k"
KERNEL=="uinput",  NAME="%k", SYMLINK="misc/%k"

# alsa devices
KERNEL=="controlC[0-9]*", NAME="%k",	SYMLINK="snd/%k"
KERNEL=="hw[CD0-9]*",     NAME="%k",	SYMLINK="snd/%k"
KERNEL=="pcm[CD0-9cp]*",  NAME="%k",	SYMLINK="snd/%k"
KERNEL=="midi[CD0-9]*",   NAME="%k",	SYMLINK="snd/%k"
KERNEL=="timer",          NAME="%k",	SYMLINK="snd/%k"
KERNEL=="seq",            NAME="%k",	SYMLINK="snd/%k"

# oss devices
KERNEL=="audio*",     NAME="%k", SYMLINK="sound/%k",   MODE="0666"
KERNEL=="dmmidi",     NAME="%k", SYMLINK="sound/%k",   MODE="0666"
KERNEL=="dsp*",       NAME="%k", SYMLINK="sound/%k",   MODE="0666"
KERNEL=="midi*",      NAME="%k", SYMLINK="sound/%k",   MODE="0666"
KERNEL=="mixer*",     NAME="%k", SYMLINK="sound/%k",   MODE="0666"
KERNEL=="sequencer*", NAME="%k", SYMLINK="sound/%k",   MODE="0666"

# input devices
KERNEL=="mice",   NAME="%k",	SYMLINK="input/%k"
KERNEL=="mouse*", NAME="%k",	SYMLINK="input/%k"
KERNEL=="event*", NAME="%k",	SYMLINK="input/%k"
KERNEL=="js*",    NAME="%k",	SYMLINK="input/%k"
KERNEL=="ts*",	 NAME="%k",	SYMLINK="input/%k"

# netlink devices
KERNEL=="route",		NAME="netlink/%k"
KERNEL=="skip",		NAME="netlink/%k"
KERNEL=="usersock",	NAME="netlink/%k"
KERNEL=="fwmonitor",	NAME="netlink/%k"
KERNEL=="tcpdiag",	NAME="netlink/%k"
KERNEL=="nflog",		NAME="netlink/%k"
KERNEL=="xfrm",		NAME="netlink/%k"
KERNEL=="arpd",		NAME="netlink/%k"
KERNEL=="route6",	NAME="netlink/%k"
KERNEL=="ip6_fw",	NAME="netlink/%k"
KERNEL=="dnrtmsg",	NAME="netlink/%k"
KERNEL=="tap*",		NAME="netlink/%k"

# CAPI devices
KERNEL=="capi",		NAME="capi20", SYMLINK="isdn/capi20"
KERNEL=="capi*",		NAME="capi/%n"

# Network devices
KERNEL=="tun",          NAME="net/%k"

# raw devices
KERNEL=="raw[0-9]*",    NAME="%k",	SYMLINK="raw/%k" 

#Zap/Dahdi Devices
KERNEL=="dahdictl",       NAME="%k",		SYMLINK="dahdi/ctl"
KERNEL=="dahditranscode", NAME="%k",		SYMLINK="dahdi/transcode"
KERNEL=="dahditimer",     NAME="%k",		SYMLINK="dahdi/timer"
KERNEL=="dahdichannel",   NAME="%k",		SYMLINK="dahdi/channel"
KERNEL=="dahdipseudo",    NAME="%k",		SYMLINK="dahdi/pseudo"
KERNEL=="dahdi[0-9]*",    NAME="%k",		SYMLINK="dahdi/%n"

KERNEL=="zapctl",         NAME="%k",		SYMLINK="zap/ctl"
KERNEL=="zaptranscode",   NAME="%k",		SYMLINK="zap/transcode"
KERNEL=="zaptimer",	  NAME="%k",		SYMLINK="zap/timer"
KERNEL=="zapchannel",     NAME="%k",		SYMLINK="zap/channel"
KERNEL=="zappseudo",      NAME="%k",		SYMLINK="zap/pseudo"
KERNEL=="zap[0-9]*",      NAME="%k",		SYMLINK="zap/%n"

#Mem
KERNEL=="null",       NAME="%k",                MODE="0666"
KERNEL=="full",       NAME="%k",                MODE="0666"

#Random
KERNEL=="random",     NAME="%k",   MODE="0666"
KERNEL=="urandom",    NAME="%k",   MODE="0666"

# Load firmware into the Xorcom Astribank device:
BUS!="usb", ACTION!="add", GOTO="xpp_usb_add_end"
SYSFS{idVendor}=="e4e4", SYSFS{idProduct}=="11[345][01]", RUN+="/etc/hotplug/usb/xpp_fxloader udev $sysfs{idVendor}/$sysfs{idProduct}/$sysfs{bcdDevice}"
LABEL="xpp_usb_add_end"
